# services/socat/run
#!/bin/sh

exec socat TCP-LISTEN:5432,fork,reuseaddr "TCP-CONNECT:$DATABASE_HOST:$DATABASE_PORT,connect-timeout=2"
