# services/pgbouncer/run
#!/bin/sh

exec su pgbouncer -c 'exec /usr/bin/pgbouncer /conf/pgbouncer.conf'
